# mail/postfix/simple/conf/supervisor/opendkim.ini
[program:opendkim]
command=/usr/sbin/opendkim -flx /etc/opendkim/opendkim.conf
priority=3
startsecs=0
