# mail/postfix/simple/conf/supervisor/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=2
startsecs=0
